# roles/elasticsearch_curator/templates/elasticsearch-curator.cron.d
0 2 * * * * root curator --host {{ hostvars[groups["elasticsearch"][0]]["ansible_fqdn"] }} {{ elasticsearch_curator_args }}
