# roles/ganglia_monitor/templates/device-metrics.cron.d
5-55/10 * * * * root php /usr/local/sbin/device-metrics.php 1 1
